# aimedbot.service
[Unit]
Description=BlackBox GPT — Universal AI Assistant Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/AIMedChatBot
ExecStart=/home/ubuntu/AIMedChatBot/venv/bin/python -m bot.main
Restart=on-failure
RestartSec=5
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
