# systemd/blackboxgpt@.service
[Unit]
Description=BlackBoxGPT Telegram Bot (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/BlackBoxGPT_bot
EnvironmentFile=/home/%i/BlackBoxGPT_bot/.env
ExecStart=/home/%i/BlackBoxGPT_bot/venv/bin/python -m bot.main
Restart=always
RestartSec=3
TimeoutStopSec=10
KillSignal=SIGINT
SyslogIdentifier=blackboxgpt
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
