# deploy/aimedbot.service
[Unit]
Description=BlackBoxGPT Telegram Bot
After=network-online.target

[Service]
Type=simple
User=your-linux-user
WorkingDirectory=/home/your-linux-user/BlackBoxGPTBot
EnvironmentFile=/home/your-linux-user/BlackBoxGPTBot/.env
ExecStart=/home/your-linux-user/BlackBoxGPTBot/venv/bin/python -m bot.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
