# deploy/aimedbot.service
[Unit]
Description=AIMedChatBot / BlackBox GPT Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
# Change these paths and user to your real ones
User=ubuntu
WorkingDirectory=/home/ubuntu/AIMedChatBot
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/ubuntu/AIMedChatBot/venv/bin/python -m bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
